# tests/myproduct/supervisor/myproduct.conf
[program:myproduct-api-8001]
command=/usr/local/bin/apiphant /opt/myproduct 127.0.0.1:8001
directory=/opt/myproduct
user=myproduct
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/opt/myproduct/log/supervisor-api-8001.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

[program:myproduct-background]
command=/usr/local/bin/apiphant-background /opt/myproduct
directory=/opt/myproduct
user=myproduct
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/opt/myproduct/log/supervisor-background.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

[group:myproduct]
programs=myproduct-background,myproduct-api-8001
# Add more here: 8002, 8003, etc.
